# scripts/careerbridge-frontend.service
[Unit]
Description=CareerBridge Frontend (React + Vite Preview)
After=network-online.target careerbridge-backend.service
Wants=network-online.target careerbridge-backend.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=iiitd
WorkingDirectory=/home/iiitd/projects/FCS/frontend
EnvironmentFile=/home/iiitd/projects/FCS/.env
Environment="NVM_DIR=/home/iiitd/.nvm"
ExecStart=/bin/bash -lc 'source "$NVM_DIR/nvm.sh" && nvm use 20 >/dev/null && export VITE_API_BASE_URL="http://192.168.3.40:8010" && npm run preview -- --host 0.0.0.0 --port 4173'
Restart=always
RestartSec=5
StandardOutput=append:/home/iiitd/projects/FCS/frontend/frontend.log
StandardError=append:/home/iiitd/projects/FCS/frontend/frontend.log

[Install]
WantedBy=multi-user.target
